# packaging/debian/lib/systemd/system/riftchat.service
[Unit]
Description=RiftChat - Reticulum mesh messenger
Documentation=https://rift.pw  https://github.com/datton/TheRift/tree/main/apps/riftchat
After=network-online.target
# Soft ordering: if rnsd-agent is also installed (Rift transport sidecar),
# let it start first so it becomes the RNS share_instance master with
# enable_transport=yes. RiftChat then attaches as a client. If rnsd-agent
# isn't installed, this After= is a no-op and RiftChat boots immediately.
After=rnsd-agent.service
Wants=network-online.target

[Service]
Type=simple
User=riftchat
Group=riftchat
WorkingDirectory=/opt/riftchat
Environment="PATH=/usr/lib/rift/venv/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin"
Environment="PYTHONIOENCODING=utf-8"
Environment="LANG=C.UTF-8"
# Runs out of the shared /usr/lib/rift/venv owned by rift-shared (one rns,
# one lxmf, one set of crypto deps across every Rift app on this host).
ExecStart=/usr/lib/rift/venv/bin/python /opt/riftchat/meshchat.py \
    --host 0.0.0.0 \
    --port 8111 \
    --headless \
    --identity-file /var/lib/riftchat/identity \
    --storage-dir /var/lib/riftchat/storage \
    --reticulum-config-dir /var/lib/riftchat/reticulum

# Readiness probe — systemd considers the unit "active" only once the UI
# answers on :8111. Retries for up to 30s with 1s intervals. Failure here
# triggers the Restart=on-failure policy below.
ExecStartPost=/bin/bash -c 'for i in $(seq 1 30); do curl -fsS -o /dev/null --max-time 2 http://127.0.0.1:8111/ && exit 0; sleep 1; done; echo "[riftchat] UI did not become reachable on :8111 within 30s" >&2; exit 1'
TimeoutStartSec=60

# Restart policy
Restart=on-failure
RestartSec=10

# Security hardening (loosened where it would block Reticulum's transport detection)
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=full
ProtectHome=true
ReadWritePaths=/var/lib/riftchat
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true

# Allow access to USB serial devices (RNode / LoRa)
SupplementaryGroups=dialout
DeviceAllow=char-ttyUSB rw
DeviceAllow=char-ttyACM rw

[Install]
WantedBy=multi-user.target
